# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:xvfb]
command=Xvfb :0 -screen 0 %(ENV_VNC_RESOLUTION)s -ac +extension GLX +render -noreset
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:openbox]
command=openbox
autorestart=true
priority=200
environment=DISPLAY=":0"
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.err

[program:x11vnc]
command=x11vnc -display :0 -rfbport 5900 -rfbauth /home/agent/.vnc/passwd -forever -shared -noxdamage -xkb -noxrecord -noxfixes -repeat
autorestart=true
priority=300
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err
startsecs=3
startretries=10

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5900
autorestart=true
priority=400
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock
autorestart=true
priority=50
user=root
stdout_logfile=/var/log/dockerd.log
stderr_logfile=/var/log/dockerd_err.log

[program:daemon]
command=python3 /opt/daemon/daemon.py
directory=/opt/daemon
user=agent
autorestart=true
priority=500
environment=DISPLAY=":0",AGENT_ID="%(ENV_AGENT_ID)s",AGENT_NAME="%(ENV_AGENT_NAME)s",REDIS_URL="%(ENV_REDIS_URL)s",AGENT_NETWORK_URL="%(ENV_AGENT_NETWORK_URL)s",PYTHONPATH="/opt:/opt/daemon:/opt/credentials"
stdout_logfile=/var/log/supervisor/daemon.log
stderr_logfile=/var/log/supervisor/daemon.err
startsecs=5
startretries=5
